# install/systemd/trashd.service
[Unit]
Description=trashd filesystem deletion monitor
Documentation=https://github.com/faratech/trashd
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/lib/trashd/trashd
Restart=always
RestartSec=5

# Capabilities
AmbientCapabilities=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
NoNewPrivileges=true

# Sandboxing
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
MemoryMax=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trashd

[Install]
WantedBy=default.target
